# systemd/picam-network.service
[Unit]
Description=PiCam Network Manager
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/pi3b-
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/admin/pi3b-/.venv/bin/python /home/admin/pi3b-/firmware/domain/network.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
